# openwrt/builds/builds/ar150/files/etc/uci-defaults/basic-lan.sh
../../../../../uci-defaults/basic-lan.sh